# config/crontab/crontab
0	*	*	*	*	root	bash /home/www/bitwarden.test/www/bin/console.sh dump >/dev/null 2>&1